# packages/mesos/mesos-agent.initd
#!/sbin/runscript

name=mesos-master
daemon=/usr/sbin/$name

#depend() {
	#need net
	#after firewall
#}

start() {
	ebegin "Starting ${name}"
    # Build --env args with variables from confd
    local opts=""
    for i in $(set |grep ^MESOS_) ; do 
      opts="$opts --env $i"
    done
		start-stop-daemon --start --quiet \
			--pidfile /var/run/${name}.pid \
      $opts --exec ${daemon}
	eend $?
}

stop() {
	ebegin "Stopping ${name}"
		start-stop-daemon --stop --quiet \
			--pidfile /var/run/$name.pid \
			--exec ${daemon}
	eend $?
}

