# pkg/sshd/etc/profile.d/mesanine.sh
#!/bin/sh
TERM=xterm-color
PS1='\e[33;1m\u@\h: \e[31m\W\e[0m\$ '
alias ll='ls -alh'
